# deploy/broker-2/glowup-ble-sensor.service
# GlowUp BLE Sensor Daemon — service-pattern producer on broker-2
#
# Runs the BLE sensor daemon (ble/sensor.py) as a service on broker-2,
# the host that owns the BT radio and the ONVIS sensor pairings.
# Publishes sensor signals CROSS-HOST directly to the hub mosquitto
# at GLB_HUB_BROKER (default 10.0.0.214).  The hub does not subscribe
# to broker-2 — that direction is the failure mode that motivated
# the service-pattern pivot on 2026-04-15.  See:
#
#   - docs/35-service-vs-adapter.md (the architectural rule)
#   - docs/28-ble-sensors.md        (the chapter for this service)
#   - feedback_read_the_producer_first.md
#   - feedback_service_vs_adapter_rule.md
#   - feedback_multi_topic_config_deletion.md
#
# Topics published (all qos=0, retain=False — never retain sensor data):
#
#   glowup/signals/{label}:motion       → "1" / "0"
#   glowup/signals/{label}:temperature  → Celsius float
#   glowup/signals/{label}:humidity     → percent float
#   glowup/ble/status/{label}           → JSON status blob
#
# The hub's _on_remote_signal callback consumes glowup/signals/#;
# a separate hub-side callback consumes glowup/ble/status/# for
# the /api/ble/sensors diagnostic endpoint.
#
# Install / update on broker-2:
#
#   scp deploy/broker-2/glowup-ble-sensor.service \
#       a@10.0.0.123:/tmp/glowup-ble-sensor.service
#   ssh a@10.0.0.123 \
#     'sudo install -o root -g root -m 0644 \
#        /tmp/glowup-ble-sensor.service \
#        /etc/systemd/system/glowup-ble-sensor.service && \
#      sudo systemctl daemon-reload && \
#      sudo systemctl restart glowup-ble-sensor && \
#      sudo systemctl status glowup-ble-sensor --no-pager'
#
# Requires on broker-2:
#   - Bluetooth powered on (sudo rfkill unblock bluetooth)
#   - /etc/glowup/ble_pairing.json with paired device(s)
#   - bleak, cryptography, paho-mqtt installed in the system python
#     (or whatever interpreter ExecStart names below)
#   - The ble/ python package present at /opt/glowup-sensors/ble/
#   - Network reachability to 10.0.0.214:1883 (hub mosquitto)
#
# Rollback:
#   The previous version of this unit passed `--broker localhost`
#   and depended on a now-deleted mosquitto bridge to forward
#   glowup/ble/# from broker-2 to the hub.  Do NOT roll back to
#   that arrangement — the bridge does not exist anymore and the
#   adapter-side consumer (adapters/ble_adapter.py) has been
#   deleted from the hub.  If a true rollback is ever needed, the
#   correct prior architecture is the in-process BleAdapter on
#   the hub, not the broken bridge.

[Unit]
Description=GlowUp BLE sensor daemon (HAP-BLE → cross-host signal publish)
After=network-online.target bluetooth.target
Wants=network-online.target bluetooth.target

[Service]
Type=simple
User=a
Group=a
WorkingDirectory=/opt/glowup-sensors

# Hub mosquitto address — the systemd unit owns this (per the
# installer-owns-config rule).  Change here if the hub IP ever
# moves and re-install the unit.
Environment=GLB_HUB_BROKER=10.0.0.214
Environment=GLB_HUB_PORT=1883

ExecStartPre=-/usr/sbin/rfkill unblock bluetooth
ExecStart=/usr/bin/python3 -m ble.sensor --config /etc/glowup/ble_pairing.json
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=glowup-ble-sensor

[Install]
WantedBy=multi-user.target
