# deploy/broker-2/glowup-ble-sensor.service
# GlowUp BLE Sensor Daemon — broker-2 (relay node)
#
# Same daemon as the primary Pi but publishes to localhost mosquitto,
# which bridges messages upstream to the primary Pi.
#
# Install:
#   sudo cp glowup-ble-sensor.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable glowup-ble-sensor
#   sudo systemctl start glowup-ble-sensor
#
# Requires:
#   - Bluetooth powered on (sudo rfkill unblock bluetooth)
#   - ble_pairing.json with paired device(s)
#   - MQTT broker running locally (mosquitto with bridge config)
#   - bleak, cryptography, paho-mqtt pip packages

[Unit]
Description=GlowUp BLE Sensor Daemon (broker-2 relay)
After=bluetooth.target mosquitto.service network-online.target
Wants=bluetooth.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/a/lifx
ExecStartPre=-/usr/sbin/rfkill unblock bluetooth
ExecStart=/home/a/venv/bin/python3 -m ble.sensor --config /home/a/lifx/ble_pairing.json --broker localhost
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
