# glowup-agent.service
[Unit]
Description=GlowUp Worker Agent
After=network-online.target mosquitto.service
Wants=network-online.target

[Service]
Type=simple
User=a
WorkingDirectory=/home/a/lifx
ExecStart=/usr/bin/python3 -m distributed.worker_agent /etc/glowup/agent.json
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=glowup-agent

# Environment
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
