# services/com.glowup.macos-thermal.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  GlowUp macOS Thermal Sensor — LaunchAgent for Daedalus (Mac Studio).

  Publishes CPU temp + load + uptime to MQTT every 30 s so Daedalus
  appears in the GlowUp thermal fleet dashboard alongside the Pis.

  This is a LaunchAgent (per-user), not a LaunchDaemon — `smctemp`
  reads the SMC without root, and we want the daemon to inherit the
  user's brew PATH.

  Install:
    brew install smctemp
    cp com.glowup.macos-thermal.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.glowup.macos-thermal.plist

  Remove:
    launchctl unload ~/Library/LaunchAgents/com.glowup.macos-thermal.plist
    rm ~/Library/LaunchAgents/com.glowup.macos-thermal.plist

  Logs:
    tail -f ~/Library/Logs/glowup-macos-thermal.log
    tail -f ~/Library/Logs/glowup-macos-thermal.err

  Requires:
    ~/macos_thermal.conf (copy from contrib/sensors/macos_thermal.conf.example)
    ~/venv with paho-mqtt installed
    ~/lifx/ synced from staging
    smctemp on PATH (brew installs to /opt/homebrew/bin on Apple Silicon)
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.glowup.macos-thermal</string>

    <key>WorkingDirectory</key>
    <string>/Users/perrykivolowitz/lifx</string>

    <!-- Run as a module (-m), not as a script.  Running as a script
         (`python contrib/sensors/macos_thermal_sensor.py`) puts the
         script's directory on sys.path[0], so `from glowup_site
         import site` fails because glowup_site.py lives at the lifx
         repo root, not under contrib/sensors/.  The -m form sets
         sys.path[0] to cwd (the WorkingDirectory above), which DOES
         contain glowup_site.py.  The same path-leak issue is why
         this whole plist arguably belongs in glowup-infra and gets
         templated per host; until that refactor lands, this one
         line keeps the LaunchAgent working. -->
    <key>ProgramArguments</key>
    <array>
        <string>/Users/perrykivolowitz/venv/bin/python</string>
        <string>-m</string>
        <string>contrib.sensors.macos_thermal_sensor</string>
        <string>--config</string>
        <string>/Users/perrykivolowitz/macos_thermal.conf</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <!-- /usr/sbin is required: macmon shells out to powermetrics
             (which lives at /usr/sbin/powermetrics) and silently fails
             with "No such file or directory" if it can't find it. -->
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/perrykivolowitz</string>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>/Users/perrykivolowitz/Library/Logs/glowup-macos-thermal.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/perrykivolowitz/Library/Logs/glowup-macos-thermal.err</string>
</dict>
</plist>
