# services/kiosk-health.service
[Unit]
# GlowUp kiosk health sensor — fleet-portable.
#
# Observes a single systemd kiosk service and publishes normalized
# health records to the local MQTT broker on:
#
#     glowup/hardware/kiosk/<node_id>       (retained, every interval)
#     glowup/node/<node_id>/kiosk/status    (retained, online/offline, LWT)
#
# Source: contrib/sensors/kiosk_health_sensor.py in the lifx repo.
# Deploy: /opt/glowup-sensors/kiosk_health_sensor.py
#
# This unit is USER scope because the kiosk we observe is user scope
# on mbclock (pygame needs $DISPLAY, which is a user session thing).
# On a kiosk host where the kiosk is a system unit, copy this file to
# /etc/systemd/system/, drop the --user flag from the ExecStart, and
# change --scope to system.
#
# Runtime dep: python3-paho-mqtt (apt or pip inside the kiosk venv).
#
# Recreate on a kiosk host:
#   mkdir -p ~/.config/systemd/user
#   cp this-file ~/.config/systemd/user/kiosk-health.service
#   systemctl --user daemon-reload && systemctl --user enable --now kiosk-health
#   loginctl enable-linger mortimer.snerd  # so the sensor runs even w/o login
Description=GlowUp kiosk health sensor
Documentation=file:///opt/glowup-sensors/kiosk_health_sensor.py
After=default.target

[Service]
Type=simple
ExecStart=/home/a/venv/bin/python /opt/glowup-sensors/kiosk_health_sensor.py \
    --broker 127.0.0.1 --service kiosk.service --scope user
Restart=always
RestartSec=10s
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
