# deploy/glowup-server.service
# Copyright (c) 2026 Perry Kivolowitz. All rights reserved.
# Licensed under the MIT License. See LICENSE file in the project root.

[Unit]
Description=GlowUp REST API Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=a
ExecStart=/usr/bin/python3 /home/a/lifx/server.py /etc/glowup/server.json
WorkingDirectory=/home/a/lifx
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=glowup-server

# Hardening — server needs write access to state.db and the repo tree.
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictSUIDSGID=yes
LockPersonality=yes
ReadWritePaths=/etc/glowup /home/a/lifx

[Install]
WantedBy=multi-user.target
