# tools/remote_hid/deploy/glowup-remote-hid.service
# glowup-remote-hid daemon — virtual HID device + TCP listener.
#
# Purpose: a Pi with an HDMI monitor but no keyboard/mouse becomes
# remote-controllable from any machine running tools/remote_hid/client.py.
# Useful for kiosk hosts (mbclock, broker-2 during diagnostics) and for
# anyone temporarily plugging a spare monitor into a Pi.
#
# Install on a Pi:
#
#   scp tools/remote_hid/*.py mortimer.snerd@pi:/tmp/remote_hid/
#   ssh mortimer.snerd@pi 'sudo mkdir -p /opt/glowup-remote-hid && \
#       sudo install -o root -g root -m 0644 /tmp/remote_hid/*.py \
#           /opt/glowup-remote-hid/'
#   scp tools/remote_hid/deploy/99-uinput.rules mortimer.snerd@pi:/tmp/
#   ssh mortimer.snerd@pi 'sudo install -o root -g root -m 0644 \
#       /tmp/99-uinput.rules /etc/udev/rules.d/ && \
#       sudo udevadm control --reload-rules && sudo udevadm trigger'
#   scp tools/remote_hid/remote_hid.json.example mortimer.snerd@pi:/tmp/
#   ssh mortimer.snerd@pi 'sudo install -o mortimer.snerd -g mortimer.snerd -m 0600 /tmp/remote_hid.json.example \
#       /etc/glowup/remote_hid.json'
#   ssh mortimer.snerd@pi 'sudo /opt/glowup-remote-hid/../venv/bin/pip install evdev'
#   scp tools/remote_hid/deploy/glowup-remote-hid.service mortimer.snerd@pi:/tmp/
#   ssh mortimer.snerd@pi 'sudo install -o root -g root -m 0644 \
#       /tmp/glowup-remote-hid.service /etc/systemd/system/ && \
#       sudo systemctl daemon-reload && \
#       sudo systemctl enable --now glowup-remote-hid'
#
# Requires on the target Pi:
#   - /dev/uinput present (CONFIG_INPUT_UINPUT built into the Pi OS kernel)
#   - User `a` in group `input` (udev rule /etc/udev/rules.d/99-uinput.rules
#     makes /dev/uinput mode 0660 group input; we also gid-assign at runtime)
#   - python3-evdev in the venv at /opt/glowup-remote-hid/venv

[Unit]
Description=GlowUp Remote HID — virtual mouse/keyboard over TCP
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=a
Group=a
SupplementaryGroups=input
WorkingDirectory=/opt/glowup-remote-hid

# PYTHONPATH lets us invoke as `python -m remote_hid.server` without
# packaging; /opt/glowup-remote-hid/ *is* the package root (symlinked
# or copied as remote_hid/).
Environment=PYTHONPATH=/opt/glowup-remote-hid
ExecStart=/opt/glowup-remote-hid/venv/bin/python -m remote_hid.server \
    --config /etc/glowup/remote_hid.json
Restart=on-failure
RestartSec=5

StandardOutput=journal
StandardError=journal
SyslogIdentifier=glowup-remote-hid

# Hardening — daemon only needs /dev/uinput and the config file.
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictSUIDSGID=yes
LockPersonality=yes
# /dev/uinput is a character device; needs DeviceAllow to coexist
# with ProtectSystem=strict.
DeviceAllow=/dev/uinput rw
DevicePolicy=closed

[Install]
WantedBy=multi-user.target
